# Laboratorio 3 - Boot y Servicios/consulta.service
[Unit]
Description=Asks for managed packages.

[Service]
Type=simple
ExecStart=/bin/bash /usr/sbin/consultad

[Install]
WantedBy=multi-user.target
